# moinho/comun/etc/uci-defaults/96-minido-info
#!/bin/sh
### Workaround dnsmasq procd init.d bug
sed -i '\|/thisnode.info/|{p;s|/thisnode.info/|/meuno.info/|}' /etc/config/dhcp
