# anisacate.libre.org.ar/comun/etc/uci-defaults/96-minido-info
#!/bin/sh
### Workaround dnsmasq procd init.d bug
sed -i '\|/thisnode.info/|{p;s|/thisnode.info/|/minodo.info/|}' /etc/config/dhcp
