# openNET.io/1144-W2PA-LIME-XXXX_ON32/etc/uci-defaults/93_dns-alias-for-anygw
#!/bin/sh

### Uncomment to create additional alias equivalent to thisnode.info
### (useful for adapting to local language)
#sed -i '\|^address=/anygw/|{p;s|^address=/anygw/|address=/tplinklogin.net/|}' /etc/dnsmasq.d/*
#sed -i '\|^address=/anygw/|{p;s|^address=/anygw/|address=/roteador.info/|}' /etc/dnsmasq.d/*
