# jardibotanic/nodes/etc/uci-defaults/zz_jardibotanic_disable_dhcp_server
#!/bin/sh

rm -f /etc/smonit/dnsmasq
/etc/init.d/dnsmasq disable
/etc/init.d/odhcpd disable
