# NuestraRed.org/hacklime/etc/uci-defaults/zfirewall
#!/bin/sh
echo "#Proxy redirect" >> /etc/firewall.user
