# PromirimLivre/comum/etc/uci-defaults/97_add-localserver-to-hosts
#!/bin/sh

### Adds local server do /etc/hosts

echo "10.%N.128.10 biblioteca" >> /etc/hosts
