# openNET.io/tl-wr841-v11/etc/uci-defaults/91-change-sshd-port
#!/bin/sh
uci set dropbear.@dropbear[0].Port='42022'
uci commit dropbear
