# sbmesh/common/etc/uci-defaults/zsbmesh_banner
#!/bin/sh
cat /etc/sbmesh/sbmesh.banner >> /etc/banner
exit 0
