# sbmesh/common/etc/uci-defaults/zsystem
#!/bin/sh
uci set system.@system[0].timezone='PST8PDT,M3.2.0,M11.1.0'
uci set system.@system[0].zonename='America/Los Angeles'
uci commit system
