# calafou/outdoor_gateway/root/etc/hotplug.d/iface/50-enforce_100M_wan
[ ifup = "$ACTION" ] && {
	[ -n "$DEVICE" ] && {
		ethtool -s "$DEVICE" autoneg off speed 100 duplex full
	}
}

