# calafou/outdoor_gateway/root/etc/hotplug.d/iface/50-enforce_100M_wan
[ ifup = "$ACTION" ] && {
	[ wan = "$INTERFACE" -a -n "$DEVICE"] && {
		ethtool -s "$DEVICE" autoneg off speed 100 duplex full
	}
}

