# quintanalibre.org.ar/comun/etc/uci-defaults/91-set_alfred
#!/bin/sh
uci set alfred.alfred.mode='slave'
uci set alfred.alfred.start_vis='0'
uci commit alfred
