# quintanalibre.org.ar/comun/etc/uci-defaults/91-set_system_opts
#!/bin/sh
uci set lime.system.timezone='ART3'
uci set lime.system.log_ip='10.5.1.4'
