# quintanalibre.org.ar/comun/etc/uci-defaults/93-reboot_on_deaf_phys
#!/bin/sh
mkdir -p /etc/crontabs ; touch /etc/crontabs/root
sed "\| \* \* \* .*reset_deaf_phys.*|d" -i /etc/crontabs/root
echo "*/20 * * * * /usr/sbin/reset_deaf_phys.sh" >> /etc/crontabs/root
