# quintanalibre.org.ar/comun/etc/uci-defaults/95-copy-lucy-index
#!/bin/sh
cp /www/luci-index.html /www/index.html

