# quintanalibre.org.ar/comun/etc/uci-defaults/add_altermundi_net.sh
#!/bin/sh
echo '
config tunOut "altermundi"
        option tunOut "altermundi"
        option network "138.255.88.0/22"' >> /etc/config/bmx6
