# antennine.org/an-metrics/root/etc/uci-defaults/90_prometheus
#!/bin/sh
unique_append() {
  [ ! -f "$2" ] && mkdir -p "$(dirname "$2")" && touch "$2";
  grep -qF "$1" "$2" || echo "$1" >> "$2"; 
}
unique_append '*/5 * * * * prometheus-push' /etc/crontabs/root
