# antennine.org/an-outdoor-8-64/root/etc/uci-defaults/90_wifi-modes
#!/bin/sh
uci set lime-node.5ghz=lime-wifi-band
uci set lime-node.5ghz.modes='ieee80211s'
uci add_list lime-node.5ghz.modes='ap'
uci commit lime-node
