# klimacamp/common/etc/uci-defaults/91-disable_passwd_login
#!/bin/sh
uci set dropbear.@dropbear[0].RootPasswordAuth='off'
uci commit dropbear

