# fuma-aonline/fumacaonline-1607/etc/uci-defaults/97_add-localserver-to-hosts
#!/bin/sh

### Adds local server do /etc/hosts

echo "10.7.128.10 www" >> /etc/hosts
