# RedINC-Comunitaria/hacklime/etc/uci-defaults/94-dnsmasq_cache_size
#!/bin/sh
uci set dhcp.@dnsmasq[0].cachesize='3000'
uci commit dhcp

