# jardibotanic/BCNmtjPlPolivaCaixa/etc/uci-defaults/zz_jardibotanic_disable_automatic_internet_sharing
#!/bin/sh

/etc/init.d/watchping disable
