# kollserola/talaia/etc/uci-defaults/zz-use-odhcpd
#!/bin/sh
uci set dhcp.odhcpd.maindhcp='1'
uci commit dhcp
