# libremesh/encrypt-11s/etc/uci-defaults/80-set-11s-encryption
#!/bin/sh
uci set lime-defaults.wifi.ieee80211s_encryption='psk2+aes'
uci set lime-defaults.wifi.ieee80211s_key 'libremesh'
