# jardibotanic/mArcFlex/etc/uci-defaults/zz_jardibotanic_enable_power_over_ethernet_pass_through
#!/bin/sh

/etc/init.d/poe-passthrough enable
