# jardibotanic/mArcFlex/etc/uci-defaults/zz_jardibotanic_enable_power_over_ethernet_pass_through
#!/bin/sh

echo "@reboot echo 1 > /sys/class/gpio/poe-passthrough/value" >> /etc/crontabs/root

