# quintanalibre.org.ar/librerouter/etc/uci-defaults/89-radio2_channel_36
#!/bin/sh

# lime-config doesn't yet support a list of channels
# workaround assigning an arbitrary channel to second 5ghz radio
uci set lime.radio2.channel_5ghz=36
uci commit lime
