# quintanalibre.org.ar/librerouter/etc/uci-defaults/90-add_altermundi_net
#!/bin/sh
uci set bmx6.altermundiV4="tunOut"
uci set bmx6.altermundiV4.tunOut="altermundiV4"
uci set bmx6.altermundiV4.network="138.255.88.0/22"
uci commit bmx6
