# canmasdeu/bridge/etc/uci-defaults/disable_dhcp
#!/bin/sh
/etc/init.d/dnsmasq stop
sleep 1
/etc/init.d/dnsmasq disable
chmod 000 /etc/smonit/dnsmasq
