# openNET.io/mir3g/etc/uci-defaults/91-allow-gatewayports
#!/bin/sh
uci set dropbear.@dropbear[0].GatewayPorts='on'
uci commit dropbear
