# labolsa.libre.org.ar/comun/etc/uci-defaults/91-set_system_opts
#!/bin/sh
uci set system.@system[0].timezone='ART3'
uci set system.@system[0].log_ip='10.5.1.4'
uci commit system
