# NonoLibre/comun/etc/uci-defaults/91-set_system_opts
#!/bin/sh
uci set system.@system[0].timezone='ART3'
uci commit system
