# NonoLibre/comun/etc/uci-defaults/94-dnsmasq_set_dhcp_range
#!/bin/sh
### Change dhcp range
sed -i 's,65536,32765,' /etc/config/dhcp

