# NonoLibre/comun/etc/uci-defaults/97-server-local
#!/bin/sh
uci del_list dhcp.@dnsmasq[0].address='/server.nono.libre.org.ar/10.157.128.1'
uci add_list dhcp.@dnsmasq[0].address='/server.nono.libre.org.ar/10.157.128.1'

