# massmesh/common/etc/uci-defaults/zmassmesh_banner
#!/bin/sh
cat /etc/massmesh/mamesh.banner >> /etc/banner
exit 0
