# massmesh/common/etc/uci-defaults/zsystem
#!/bin/sh
uci set system.@system[0].timezone='EST5EDT,M3.2.0,M11.1.0'
uci set system.@system[0].zonename='America/New York'
uci commit system
